# files/nextcloud.socket
[Unit]
Description=nextcloud uwsgi socket

[Socket]
ListenStream=/run/nextcloud.sock
SocketMode=0666

[Install]
WantedBy=sockets.target
