# files/nextcloud-first-run.service.in
[Unit]
Description=nextcloud initialize config file
ConditionPathExists=!/var/lib/nextcloud/config/portable.config.php

[Service]
Type=oneshot
DynamicUser=yes
User=nextcloud
StateDirectory=nextcloud
TemporaryFileSystem=/var
ExecStart=@coreutils@/bin/mkdir -p /var/lib/nextcloud/config /var/lib/nextcloud/apps
ExecStart=@coreutils@/bin/cp @portableConfig@ /var/lib/nextcloud/config/portable.config.php
ExecStart=@coreutils@/bin/touch /var/lib/nextcloud/config/CAN_INSTALL
