# files/nextcloud-first-run.service.in
[Unit]
Description=nextcloud initialize config file
ConditionPathExists=!/var/lib/nextcloud/config/config.php

[Service]
Type=oneshot
User=nextcloud
DynamicUser=yes
StateDirectory=nextcloud
ExecStart=@coreutils@/bin/mkdir -p /var/lib/nextcloud/config /var/lib/nextcloud/apps
ExecStart=@coreutils@/bin/cp @nextcloudConfig@ /var/lib/nextcloud/config/config.php
ExecStart=@coreutils@/bin/touch /var/lib/nextcloud/config/CAN_INSTALL
