# files/nextcloud.service.in
[Unit]
Description=nextcloud uwsgi service
Wants=nextcloud-first-run.service
Requires=nextcloud.socket
After=nextcloud.socket nextcloud-first-run.service

[Service]
Type=notify
DynamicUser=yes
Environment=NEXTCLOUD_CONFIG_DIR=/var/lib/nextcloud/config/
StateDirectory=nextcloud
LogsDirectory=nextcloud
TemporaryFileSystem=/var
ExecStart=@execStart@

[Install]
WantedBy=multi-user.target
Also=nextcloud.socket
